# tools/cooja/apps/mobility/contikiprojects-code-322/edf.fr/Emma-0.2.1/EmmaSheevaPlug/80-EMMA-DISCOVERY.rules
# This file runs and stops EMMA Network Discovery
ACTION=="add", SUBSYSTEM=="usb", SYSFS{idVendor}=="03eb", RUN+="/opt/emma/core/DiscoverNetwork start"
